# aide-email.cron
#!/bin/bash

# Email the results of an AIDE database update (including the error stream)
OUTPUT=$( { /usr/sbin/aide --update; } 2>&1 )
SUBJECT="AIDE update "$(hostname)
echo "$OUTPUT" | mailx -s "$SUBJECT" root@localhost

# Use the new database
mv -f /var/lib/aide/aide.db.new.gz /var/lib/aide/aide.db.gz

